# user/acme.sh.timer
[Unit]
Description=Run acme.sh daily

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
